# hourly-time-signal.service
[Unit]
Description=Hourly Time Signal

[Install]
WantedBy=multi-user.target

[Service]
ExecStart=/usr/share/hourly-time-signal/hts.sh
